# tools/Packet.SoundModem.NinoCompare/deploy/gb7rdg-prune.timer
[Unit]
Description=Prune GB7RDG audio ring buffer every 15 min

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
